# django_ganglioside/deployment/systemd/ganglioside-django.service
[Unit]
Description=Ganglioside Django Application (Gunicorn)
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=notify
User=ganglioside
Group=ganglioside
WorkingDirectory=/var/www/ganglioside
Environment="PATH=/var/www/ganglioside/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=config.settings.production"
EnvironmentFile=/var/www/ganglioside/.env.production

ExecStart=/var/www/ganglioside/venv/bin/gunicorn \
    --config /var/www/ganglioside/gunicorn.conf.py \
    config.wsgi:application

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
